# gunicorn_freeleaps.socket
[Unit]
Description=Freeleaps socket
[Socket]
ListenStream= /home/wwwadmin/git/freeleaps/freeleaps.sock 
[Install]
WantedBy=sockets.target
